# systemd/fpp-wled-proxy.service
[Unit]
Description=FPP WLED API Proxy — HTTP Server on Port 9000
Documentation=https://github.com/bbene/fpp-WLEDProxy
After=network-online.target
Wants=network-online.target
ConditionPathExists=/home/fpp/media/plugins/fpp-WLEDProxy

[Service]
Type=simple
User=fpp
Group=fpp
WorkingDirectory=/home/fpp/media/plugins/fpp-WLEDProxy/www

# PHP built-in web server with router script
ExecStart=/usr/bin/php -S 0.0.0.0:9000 router.php

# Restart if it crashes
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fpp-wled-proxy

[Install]
WantedBy=multi-user.target
